# etc/rc.local
#!/bin/bash -e

ir-keytable -c -w /etc/rc_keymaps/harmony_mce_keyboard

exit 0
