# etc/rc.local
#!/bin/bash -e

ir-keytable -c -w /etc/rc_keymaps/imon_mce

exit 0
